# systemd/mediamtx.service
[Unit]
Description=MediaMTX RTSP Server
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mediamtx /home/pi/TheraView/config/mediamtx.yml
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
