# systemd/theraview-web.service
[Unit]
Description=TheraView Web Viewer
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /home/pi/TheraView/scripts/theraview_web.py
Restart=always
RestartSec=2
Environment=OUTDIR=/home/pi/TheraView/recordings
Environment=LOGDIR=/home/pi/TheraView/logs
Environment=WEB_PORT=8080
Environment=LIVE_VIEW_PORT=8888
Environment=RECORDER_SERVICE=theraview-camera.service

[Install]
WantedBy=multi-user.target
