# systemd/theraview-camera.service
[Unit]
Description=TheraView Camera Recording Service
After=network.target mediamtx.service
Wants=network.target mediamtx.service

[Service]
Type=simple
ExecStart=/home/pi/TheraView/scripts/theraview_recorder.sh
Restart=always
RestartSec=2
Environment=OUTDIR=/home/pi/TheraView/recordings
Environment=LOGDIR=/home/pi/TheraView/logs
Environment=CAMERA=/dev/video0
Environment=RTSP_HOST=127.0.0.1
Environment=RTSP_PORT=8554
Environment=RTSP_PATH=live

[Install]
WantedBy=multi-user.target
